# etc/systemd/system/auto-update.service
[Unit]
Description=Automatic system updater using yay

[Service]
Type=oneshot
ExecStart=/usr/local/bin/auto_update.sh
